# config/paii_poll_supervisor.ini.template
# supervisord spec to run paii poller

[program:paii]
command = /bin/python3.8 paii.paii_poll
user = pms
directory = /home/pms/src/purple_air/log
umask = 002
autostart = true
autorestart = true
